# bot.service
[Unit]
Description=Personal Assistant Telegram Bot
# Wait for full network + Tailscale MagicDNS (100.100.100.100) before start.
# Fixes startup race: get_me() -> api.telegram.org DNS lookup failing with
# "Temporary failure in name resolution" when launched before the resolver is up.
After=network-online.target tailscale-online.target
Wants=network-online.target tailscale-online.target

[Service]
User=sanel
WorkingDirectory=/home/sanel/personal-assistant-bot
ExecStart=/home/sanel/personal-assistant-bot/venv/bin/python main.py
Restart=always
RestartSec=5
Environment=HOME=/home/sanel
Environment=PYTHONWARNINGS=ignore::ResourceWarning

[Install]
WantedBy=multi-user.target
